# infra/prod/joint/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=node /app/server/dist/src/main
autostart=true
autorestart=true
redirect_stderr=true

[program:frontend]
command=serve -s /app/front/build
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
redirect_stderr=true