# packages/twenty-docker/twenty-app-dev/rootfs/etc/s6-overlay/s6-rc.d/redis/run
#!/bin/sh
exec redis-server \
  --dir /data/redis \
  --maxmemory-policy noeviction \
  --bind 127.0.0.1 \
  --protected-mode yes
