# packages/twenty-docker/twenty-app-dev/rootfs/etc/s6-overlay/s6-rc.d/init-db/up
/bin/sh /etc/s6-overlay/scripts/init-db.sh
