# packages/twenty-docker/twenty-app-dev/rootfs/etc/s6-overlay/s6-rc.d/twenty-server/run
#!/bin/sh
if [ -f /tmp/generated-app-secret ]; then
  export APP_SECRET=$(cat /tmp/generated-app-secret)
fi
cd /app/packages/twenty-server
exec yarn start:prod
