# packages/twenty-docker/twenty-dev/rootfs/etc/s6-overlay/s6-rc.d/twenty-worker/run
#!/bin/sh
cd /app/packages/twenty-server
exec node dist/src/queue-worker/queue-worker
