# packages/twenty-docker/twenty-app-dev/rootfs/etc/s6-overlay/s6-rc.d/register-crons/up
/bin/sh /etc/s6-overlay/scripts/register-crons.sh
